# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command = /usr/sbin/nginx
user = root
autostart = true

[program:php5-fpm]
command = /usr/sbin/php5-fpm -FR
user = root
autostart = true

[program:hhvm-fastcgi]
command = hhvm --mode server -vServer.Type=fastcgi -vServer.FileSocket=/var/run/hhvm/hhvm.sock
user = root
autostart = true
