# images/memcached/supervisord.conf
[supervisord]
nodaemon=true

[program:memcached]
command = /usr/sbin/memcached
user = root
autostart = true
