# images/mysql/supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
command = /usr/sbin/mysqld
user = root
autostart = true
